# config/timers/tars-exposure.timer
[Unit]
Description=T.A.R.S Public Exposure Check (daily 02:00 UTC)

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
