# config/timers/tars-memory-context.timer
[Unit]
Description=T.A.R.S Memory Context Regen (every 30 min)

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
