# config/timers/tars-health-audit.timer
[Unit]
Description=T.A.R.S Health Audit (every 6 hours, offset)

[Timer]
OnBootSec=20min
OnCalendar=*-*-* 00,06,12,18:30:00
Persistent=true

[Install]
WantedBy=timers.target
