# config/timers/tars-health-audit.timer
[Unit]
Description=T.A.R.S Health Audit (every 12 hours)

[Timer]
OnBootSec=20min
OnCalendar=*-*-* 00,12:30:00
Persistent=true

[Install]
WantedBy=timers.target
