# config/timers/tars-memory-decay.timer
[Unit]
Description=T.A.R.S Memory Decay + Stats (daily)

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
